# Shell/profile.d/in.sh
#!/bin/sh
###############################################################################
#                                                                             #
#   Filename        in.sh                                                     #
#                                                                             #
#   Description     Create a `in` alias for SysToolsLib's `In` command        #
#                                                                             #
#   Note            Using this `in` alias allows to work around the inability #
#                   to name that command `in`, as 'in' is a reserved keyword  #
#                   of the Shell language.                                    #
#                                                                             #
#   History                                                                   #
#    2023-01-09 JFL Created this script.                                      #
#                                                                             #
#                   © Copyright 2023 Jean-François Larvoire                   #
# Licensed under the Apache 2.0 license - www.apache.org/licenses/LICENSE-2.0 #
###############################################################################

alias in=In
