# Shell/profile.d/pname.sh
#!/bin/sh
###############################################################################
#                                                                             #
#   Filename        pname.sh                                                  #
#                                                                             #
#   Description     Define aliases for invoking pname and ppname.			  #
#                                                                             #
#   Properties      jEdit local buffer properties: :encoding=utf-8:tabSize=4: #
#                                                                             #
#   Note            To have the routine available in every shell, put this    #
#                   script in /etc/profile.d/.                                #
#                                                                             #
#   History                                                                   #
#    2025-04-19 JFL Created this script.                                      #
#                                                                             #
#                   © Copyright 2025 Jean-François Larvoire                   #
# Licensed under the Apache 2.0 license - www.apache.org/licenses/LICENSE-2.0 #
###############################################################################

alias pname=". pname"
alias ppname=". ppname"
